# unit-systemd/compacta-log.timer
[Unit]
Description=Compactacao log Apache Timer

[Timer]
OnCalendar=*-*-* 00:00:00
